# systemd/minecraft-backup.timer
[Unit]
Description=Hourly Backups of LanBros MC Survival Server

[Timer]
OnCalendar=00/3:00

[Install]
WantedBy=timers.target

