# systemd/minecraft.timer
[Unit]
Description=Hourly Backups of LanBros MC Survival Server

[Timer]
OnCalendar=00/3:30

[Install]
WantedBy=timers.target

